# root/etc/s6-overlay/s6-rc.d/init-endlessh-config/run
#!/usr/bin/with-contenv bash

if [ "${LOGFILE}" == "true" ]; then
    mkdir -p /config/logs/endlessh
fi

# permissions
chown -R abc:abc \
    /config \
