# root/etc/s6-overlay/s6-rc.d/init-endlessh-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

if [[ "${LOGFILE}" == "true" ]]; then
    mkdir -p /config/logs/endlessh
fi

# permissions
lsiown -R abc:abc \
    /config \
